# root/etc/services.d/duplicati/finish
#!/usr/bin/execlineb -S1

# bring down service when no password
if { s6-test ${1} -eq 42 }

s6-svscanctl -t /var/run/s6/services